# cron/prod.cron
45 2 * * * docker run --rm --name gs_orion_updateTv_prod --env-file ~/dev/www.aleonard.us-docker/orion/env/prod.env --network=gs_network --entrypoint /data/runcommand-entrypoint.sh aleonard9/orion-web:local /var/www/cgi-bin/tv /var/www/cgi-bin/tv/updatetv.php >> ~/dev/www.aleonard.us-docker/cron/log/prod/updatetv.cron.log
45 3 * * * docker run --rm --name gs_orion_updateEps_prod --env-file ~/dev/www.aleonard.us-docker/orion/env/prod.env --network=gs_network --entrypoint /data/runcommand-entrypoint.sh aleonard9/orion-web:local /var/www/cgi-bin/tv /var/www/cgi-bin/tv/updateepisodes.php >> ~/dev/www.aleonard.us-docker/cron/log/prod/updateepisodes.cron.log
45 4 * * * docker run --rm --name gs_orion_updateRecEps_prod --env-file ~/dev/www.aleonard.us-docker/orion/env/prod.env --network=gs_network --entrypoint /data/runcommand-entrypoint.sh aleonard9/orion-web:local /var/www/cgi-bin/tv /var/www/cgi-bin/tv/updaterecentepisodes.php >> ~/dev/www.aleonard.us-docker/cron/log/prod/updaterecentepisodes.cron.log
30 2 * * * docker run --rm --name gs_orion_updateTv_dev --env-file ~/dev/www.aleonard.us-docker/orion/env/dev.env --network=gs_network_dev --entrypoint /data/runcommand-entrypoint.sh aleonard9/orion-web:local /var/www/cgi-bin/tv /var/www/cgi-bin/tv/updatetv.php >> ~/dev/www.aleonard.us-docker/cron/log/dev/updatetv.cron.log
45 3 * * * docker run --rm --name gs_orion_updateEps_dev --env-file ~/dev/www.aleonard.us-docker/orion/env/dev.env --network=gs_network_dev --entrypoint /data/runcommand-entrypoint.sh aleonard9/orion-web:local /var/www/cgi-bin/tv /var/www/cgi-bin/tv/updateepisodes.php >> ~/dev/www.aleonard.us-docker/cron/log/dev/updateepisodes.cron.log
45 4 * * * docker run --rm --name gs_orion_updateRecEps_dev --env-file ~/dev/www.aleonard.us-docker/orion/env/dev.env --network=gs_network_dev --entrypoint /data/runcommand-entrypoint.sh aleonard9/orion-web:local /var/www/cgi-bin/tv /var/www/cgi-bin/tv/updaterecentepisodes.php >> ~/dev/www.aleonard.us-docker/cron/log/dev/updaterecentepisodes.cron.log
10 5 * * * ~/dev/www.aleonard.us-docker/database/backup/mysqldump.sh dev >> ~/dev/www.aleonard.us-docker/database/backup/log/dev/backup.log
15 5 * * * ~/dev/www.aleonard.us-docker/database/backup/mysqldump.sh prod >> ~/dev/www.aleonard.us-docker/database/backup/log/prod/backup.log